# ansible/playbooks/templates/systemd-custom.service.j2
[Unit]
After={{ after }}
[Service]
ExecStart={{ exec_start }}
[Install]
WantedBy={{ wanted_by }}
